# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:backend]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000
directory=/app/backend
environment=DATABASE_URL="sqlite:////app/data/sportshield.db",REDIS_URL="redis://localhost:6379/0",JWT_SECRET="%(ENV_JWT_SECRET)s",UPLOAD_DIR="/app/data/uploads",FAISS_INDEX_PATH="/app/data/faiss/index.bin",THRESHOLD="0.80"
autostart=true
autorestart=true
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend.log

[program:worker]
command=celery -A app.worker.celery_app worker --loglevel=info --concurrency=2
directory=/app/backend
environment=DATABASE_URL="sqlite:////app/data/sportshield.db",REDIS_URL="redis://localhost:6379/0",JWT_SECRET="%(ENV_JWT_SECRET)s",UPLOAD_DIR="/app/data/uploads",FAISS_INDEX_PATH="/app/data/faiss/index.bin",THRESHOLD="0.80"
autostart=true
autorestart=true
stdout_logfile=/var/log/worker.log
stderr_logfile=/var/log/worker.log

[program:beat]
command=celery -A app.worker.celery_app beat --loglevel=info
directory=/app/backend
environment=DATABASE_URL="sqlite:////app/data/sportshield.db",REDIS_URL="redis://localhost:6379/0",JWT_SECRET="%(ENV_JWT_SECRET)s",UPLOAD_DIR="/app/data/uploads",FAISS_INDEX_PATH="/app/data/faiss/index.bin",THRESHOLD="0.80"
autostart=true
autorestart=true
stdout_logfile=/var/log/beat.log
stderr_logfile=/var/log/beat.log

[program:frontend]
command=node .next/standalone/server.js
directory=/app/frontend
environment=PORT="7860",NEXT_PUBLIC_API_URL="%(ENV_SPACE_URL)s",NEXT_PUBLIC_WS_URL="%(ENV_SPACE_URL)s"
autostart=true
autorestart=true
stdout_logfile=/var/log/frontend.log
stderr_logfile=/var/log/frontend.log
